# board/common/overlay/etc/profile.d/path-local.sh
# add local to standard path
export PATH=$PATH:/usr/local/bin:/usr/local/sbin
