# board/common/overlay/etc/profile.d/profile-rc.sh
# always execute /etc/profile for new shells
export ENV=/etc/profile
